# update.cron
0/20 * * * * /var/www/podcast/update.sh
