# update.cron
*/20 * * * * /home/pi/mradio/update.sh
*/2 * * * * /home/pi/mradio/random.sh
