# update.cron
*/20 * * * * /home/pi/mradio/update.sh
*/10 * * * * /home/pi/cloudflare/updateCloudflare.sh
